# gen/services/venom-bot.service
[Unit]
Description=Venom Bot
After=network.target

[Service]
WorkingDirectory={{appRoot}}
ExecStart=/root/Library/pnpm/pnpm start
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target

