# gen/services/venom-bot.service
[Unit]
Description=Venom Bot
After=network.target

[Service]
WorkingDirectory={{appRoot}}
ExecStart=/usr/local/bin/dotenvx run --env-file=.env.local -- /root/Library/pnpm/pnpm start
ExecStop=/bin/kill -SIGINT $MAINPID
KillSignal=SIGINT
SendSIGKILL=no
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target

